# ridgeback_bringup/scripts/ros.service
# THIS IS A GENERATED FILE, NOT RECOMMENDED TO EDIT.

[Unit]
Description="bringup ros"
After=can-udp-bridge.service

[Service]
Type=simple
ExecStart=/usr/sbin/ros-start
#ExecStart=/usr/sbin/ros2-start

[Install]
WantedBy=multi-user.target
